# etc/init.d/mining
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
        need net
}

start() {
        ebegin "Starting ${PROG} in screen"
	start-stop-daemon \
	--start --user "${SCREEN_OWNER}" --name ${PROG} \
	--env HOME="${VARDIR}" --exec screen \
	--nicelevel "${NICELEVEL}" \
	--background \
	-- ${OPTS}
	eend $?

}

stop() {
	#PID=`find /home/${SCREEN_OWNER}/.screen -name *.${PROG}`
        ebegin "Stopping ${PROG}"
	start-stop-daemon --stop --user "${SCREEN_OWNER}" \
	--name ${PROG} \
	--wait 30000 \
	--progress
	eend $?
}

